# docker/systemd.target.erb
[Unit]
Description=kegbot/<%= server_name %> container
Requires=docker.service
After=docker.service

[Service]
Restart=always
ExecStart=/usr/bin/docker start -a -p 8000:<%= server_port %> kegbot/<%= server_name %>
ExecStop=/usr/bin/docker stop -t 2 kegbot/<%= server_name %>

[Install]
WantedBy=local.target
